# deploy/systemd/homecam-detect.service
[Unit]
Description=Home Camera detection worker (jetson-inference)
Requires=mediamtx.service
After=mediamtx.service homecam-server.service
Wants=homecam-server.service
# Restart-storm guard: if the worker fails 5 times within 60 s (i.e. an
# unrecoverable bug like a Python syntax error or a missing module),
# systemd stops auto-restarting and surfaces the failure rather than
# spinning the cgroup at 100 %. iter-28's `from __future__ import
# annotations` (Py 3.7+ syntax on the 3.6 host) hit ~42 restarts in
# seconds before someone noticed; this clamps that.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
# iter (camera-recovery 2026-06-20): sd_notify liveness watchdog. detect.py
# sends READY=1 once the model + camera are up (so the long TRT load doesn't
# count against WatchdogSec), then pings WATCHDOG=1 from a DEDICATED thread
# every ~20 s. The first attempt pinged from the MAIN LOOP and cycle-failed —
# the loop can block in camera.Capture() during stream startup, so the ping
# missed the window. The dedicated thread can't block, so it reliably refreshes
# the timer; it proves the PROCESS is alive (an OOM/kernel freeze stops the
# thread → systemd restarts, persisted escalation resumes). NotifyAccess=all is
# REQUIRED: run-detect.sh (main PID) forks python, so the pings come from a
# CHILD — the default NotifyAccess=main rejects them. WatchdogSec is generous
# (4.5x the ping interval). TimeoutStartSec is long for the first-ever TRT
# COMPILE before READY (cached deserialise is ~6 s).
Type=notify
NotifyAccess=all
WatchdogSec=90
TimeoutStartSec=600
User=israel
Group=israel
# `video` is required for jetson-utils' `videoSource` to open NVDEC for
# H.264 decode of the MediaMTX RTSP stream. Removing this group breaks the
# capture path with EACCES on /dev/nvhost-* nodes.
SupplementaryGroups=video
WorkingDirectory=/home/israel/HomeCameraSystem
Environment=EVENT_URL=http://127.0.0.1:8000/api/_internal/event
Environment=DETECT_SOURCE=rtsp://localhost:8554/cam
Environment=DETECT_THRESHOLD=0.55
Environment=DETECT_COOLDOWN_S=5.0
Environment=DETECT_MODEL=ssd-mobilenet-v2
ExecStart=/home/israel/HomeCameraSystem/detection/run-detect.sh
Restart=on-failure
RestartSec=5
# nvarguscamerasrc/jetson-utils touches /tmp; don't enable PrivateTmp (same
# trap mediamtx.service hit — see deploy/systemd/mediamtx.service).
#
# iter-356.63 (camera-recovery sharp edge): NoNewPrivileges INTENTIONALLY
# REMOVED. The iter-302 auto-recovery (`escalate_argus_recovery` in
# detect.py:493) shells `sudo -n systemctl restart {nvargus-daemon,
# mediamtx}` when the libargus broker wedges. With NNP=yes the kernel
# blocks setuid elevation in this unit's process group; sudo fails with
# the misleading "is /usr/bin/sudo on a file system with the 'nosuid'
# option set" error and the auto-recovery silently drops — the live
# feed stays broken until an operator SSHes in. Defense-in-depth cost
# of dropping NNP is small: `israel ALL NOPASSWD: ALL` is already on
# the host (see /etc/sudoers.d), so any process that can replace this
# unit's binary already has full root via that path. The defense bar
# is "user israel doesn't get tricked into running foreign code under
# this unit," not "this unit is sandboxed."

[Install]
WantedBy=multi-user.target
