# deploy/systemd/homecam-server.service
[Unit]
Description=Home Camera FastAPI server (Docker)
Requires=docker.service
RequiresMountsFor=/srv/homecam-media
After=docker.service srv-homecam\x2dmedia.mount
BindsTo=srv-homecam\x2dmedia.mount
# Restart-storm guard (iter-167). Pre-iter-167 this unit had no recovery
# policy at all — if `docker compose up` failed (daemon transient, image
# pull race, name collision after a manual `docker rm`), systemd marked
# the unit `failed` once and stopped trying. Asymmetric vs the peer
# units `mediamtx.service` and `homecam-detect.service` which both have
# StartLimitBurst clamps. Now matches: 3 attempts in 5 minutes, then
# surface the failure. Tighter than detect.service's 5/60 because
# `docker compose up` failures are usually structural (image build
# broken, volume mount missing) — retrying 5× rarely helps and burns
# time before the operator notices.
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/israel/HomeCameraSystem
ExecStartPre=/bin/mountpoint -q /srv/homecam-media
ExecStartPre=/usr/bin/test -f /srv/homecam-media/.homecam-media
ExecStart=/usr/bin/docker compose -f deploy/docker-compose.yml up -d
# MediaMTX delegates every publish/read decision to this API. Do not mark the
# unit active merely because Docker created the container; wait until FastAPI
# can answer before dependants launch their RTSP publishers.
ExecStartPost=/bin/sh -c 'for _ in 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15 16 17 18 19 20 21 22 23 24 25 26 27 28 29 30 31 32 33 34 35 36 37 38 39 40 41 42 43 44 45; do /usr/bin/curl --fail --silent --max-time 2 http://127.0.0.1:8000/healthz && exit 0; /bin/sleep 1; done; exit 1'
ExecStop=/usr/bin/docker compose -f deploy/docker-compose.yml down
# Type=oneshot is intentionally not paired with Restart=. JetPack's older
# systemd rejects that combination and refuses to load the entire unit.
# Cap shutdown wait. Default systemd timeout is 90 s; combined with
# `docker compose down`'s graceful-stop attempt this can stall reboots.
# 30 s is plenty for a healthy stop and bounds an ungraceful one.
TimeoutStopSec=30
User=israel
Group=israel

[Install]
WantedBy=multi-user.target
