# deploy/systemd/mediamtx.service
[Unit]
Description=MediaMTX video gateway (Home Camera System)
# iter-356.52: order/want nvargus-daemon. Without this, mediamtx's
# runOnInit gst-launch fires before nvargus-daemon has finished its
# socket setup → first attempt hits "AlreadyAllocated" on
# nvarguscamerasrc, the gst-launch crashes, mediamtx restarts, the
# next attempt races again. Boot logs showed 14 gst-launch processes
# tried and failed across a 5-minute window before stabilising. The
# After= alone establishes ordering; the Wants= ensures nvargus-daemon
# is actually started (it is socket-activated and may not be active
# at our boot time).
After=network-online.target nvargus-daemon.service
Wants=network-online.target
# iter (camera-recovery 2026-06-20, deep-research-backed): BindsTo makes the
# Argus daemon's runtime liveness a HARD precondition — mediamtx is torn down
# if nvargus-daemon stops/fails (it can't publish without Argus), and the
# detection watchdog's `systemctl restart nvargus-daemon` escalation correctly
# cycles mediamtx with it. BindsTo propagates stops on unexpected failure
# (stronger than Requires=) and implies the requirement dep (so the old
# Wants=nvargus is now redundant); After= above orders it. BindsTo propagates
# the STOP only — the watchdog still does an explicit `restart mediamtx` after
# the nvargus restart to bring the publisher back.
BindsTo=nvargus-daemon.service
# Match the homecam-detect storm guard: 5 failures inside 60 s and
# systemd gives up rather than spinning. A misconfigured mediamtx.yml
# fails fast and stays failed — better to see "service failed" in the
# UI's worker_alive=false signal than to flood the journal.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=israel
Group=israel
WorkingDirectory=/home/israel/HomeCameraSystem
ExecStart=/home/israel/HomeCameraSystem/mediamtx /home/israel/HomeCameraSystem/deploy/mediamtx.yml
Restart=on-failure
# iter (camera-recovery 2026-06-20): 5s→15s. A short delay hammers a still-
# wedged Argus on restart-on-failure (deep research: use an explicit longer
# RestartSec, 10-30s, since systemd geometric backoff via RestartSteps is not
# reliably available here). The escalating detection watchdog owns the
# heavier nvargus/reboot recovery; this just paces mediamtx's own restarts.
RestartSec=15
# nvarguscamerasrc opens the libargus IPC socket — grant the user the
# `video` group so it has permission. This is normally already the case for
# a desktop user on JetPack, but make it explicit:
SupplementaryGroups=video
# nvarguscamerasrc talks to nvargus-daemon over /tmp/argus_socket; PrivateTmp
# would hide that socket and the GStreamer pipeline would fail with
# "Connecting to nvargus-daemon failed: No such file or directory". Don't
# enable PrivateTmp here.
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
